# config/messenger-consume.ini
[program:messenger-consume]
command=php /Users/melvinkoopmans/Development/www/code-messenger/start/bin/console messenger:consume async_priority_high async --time-limit=3600
user=melvinkoopmans
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d